# deploy/systemd/atchess-protocol.service
[Unit]
Description=ATChess Protocol Service
After=network.target

[Service]
Type=simple
User=atchess
Group=www
WorkingDirectory=/srv/atchess/app

# Load environment variables from file
EnvironmentFile=/etc/atchess/protocol.env
ExecStart=/srv/atchess/app/atchess-protocol

Restart=on-failure
RestartSec=5

# Logging - using syslog for compatibility with older systemd
StandardOutput=syslog
StandardError=syslog

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/atchess/data /srv/atchess/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
# PrivateMounts requires systemd 233+

# Resource limits
LimitNOFILE=1024
# MemoryLimit is deprecated, use MemoryMax for newer systemd
# CPUQuota requires specific syntax

# Logging - append: requires systemd 240+
SyslogIdentifier=atchess-protocol
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target